# nginx-php/7.2-fpm-stretch/supervisord/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F --fpm-config=/usr/local/etc/php-fpm.conf
autorestart=true
